# systemd/ubopt-agent.service
[Unit]
Description=UbuntOptimizer Agent
Documentation=https://github.com/120git/tinyllama-x
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ubopt update --headless --json --config /etc/ubopt/ubopt.yaml
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ubopt /var/log/ubopt
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
